# schedulers/macos/com.quietpatch.catalog-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.quietpatch.catalog-update</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/quietpatch</string>
    <string>catalog-update</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>$(HOME)/.quietpatch/logs/update.out</string>
  <key>StandardErrorPath</key>
  <string>$(HOME)/.quietpatch/logs/update.err</string>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>QP_OFFLINE</key>
    <string>1</string>
  </dict>
</dict>
</plist>
